# backend/helper-scripts/com.aibacktest.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label: Unique identifier for this job -->
    <key>Label</key>
    <string>com.aibacktest.backup</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/edwardkim/Code/ai-backtest/backend && /usr/local/bin/npm run backup >> /Users/edwardkim/Backups/ai-backtest/backup.log 2>&1</string>
    </array>

    <!-- Schedule: Daily at 2:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Run even if the system was asleep at the scheduled time -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/edwardkim/Backups/ai-backtest/backup.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/edwardkim/Backups/ai-backtest/backup-error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/edwardkim/Code/ai-backtest/backend</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
